# opencrab.service
# File path: /etc/systemd/system/opencrab.service
# sudo systemctl daemon-reload
# sudo systemctl start opencrab
# sudo systemctl enable opencrab
# sudo systemctl status opencrab
[Unit]
Description=One API Service
After=network.target

[Service]
User=ubuntu  # 注意修改用户名
WorkingDirectory=/path/to/opencrab  # 注意修改路径
ExecStart=/path/to/opencrab/opencrab --port 3000 --log-dir /path/to/opencrab/logs  # 注意修改路径和端口号
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
